# flextrend-bot/scripts/systemd/flextrend-bot.service
[Unit]
Description=FlexTrend Bot (FastAPI + Web Monitor)
After=network.target

[Service]
WorkingDirectory=/home/USER/flextrend-bot
Environment=APP_HOST=0.0.0.0
ExecStart=/home/USER/flextrend-bot/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
